# ansible/roles/tothemoon/templates/tothemoon-ws.service.j2
[Unit]
Description=To The Moon Pump.fun WebSocket Worker
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
EnvironmentFile=/etc/tothemoon.env
ExecStart={{ app_dir }}/venv/bin/python -m src.workers.pumpfun_ws
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

