# ansible/roles/tothemoon/templates/tothemoon.service.j2
[Unit]
Description=To The Moon API (FastAPI + Uvicorn)
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
EnvironmentFile=/etc/tothemoon.env
ExecStart={{ app_dir }}/venv/bin/python -m uvicorn src.app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

