# scripts/systemd/tothemoon-ws.service
[Unit]
Description=To The Moon Pump.fun WebSocket Worker
After=network.target

[Service]
Type=simple
User=tothemoon
Group=tothemoon
WorkingDirectory=/srv/tothemoon
EnvironmentFile=/etc/tothemoon.env
ExecStart=/srv/tothemoon/venv/bin/python -m src.workers.pumpfun_ws
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
